# meta-adi-adsp-sc5xx/recipes-adi/gpio-boot-trace/files/gpio-boot-trace.service
[Unit]
Description=ADI GPIO boot-trace marker
After=boot-complete.target
Requires=boot-complete.target

[Service]
Type=exec
ExecStart=/usr/libexec/gpio-boot-trace

[Install]
WantedBy=multi-user.target
